# systemd/user/wallpaper-fetch.service
[Unit]
Description=Daily wallpaper fetch

[Service]
Type=oneshot
ExecStart=%h/.local/share/dotfiles/hypr/scripts/wallpaper-fetch
